# source/systemd/sotd.service
[Unit]
Description=sotd

[Service]
ExecStart=/usr/local/bin/sotd
StandardInput=tty
Type=oneshot

