# source/sotd.timer
[Unit]
Description=run sotd daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
