# source/sotd.service
[Unit]
Description=sotd

[Service]
ExecStart=/usr/bin/sotd

